# packaging/debian/voxtype.service
[Unit]
Description=Voxtype push-to-talk voice-to-text daemon
Documentation=man:voxtype(1)
PartOf=graphical-session.target
After=graphical-session.target pipewire.service pipewire-pulse.service
Wants=ydotool.service

[Service]
Type=simple
ExecStart=/usr/bin/voxtype daemon
Restart=on-failure
RestartSec=5

# Note: User must be in 'input' group for evdev access
# Before enabling this service, run: voxtype setup --download

[Install]
WantedBy=graphical-session.target
