# packaging/debian/voxtype.service
[Unit]
Description=Voxtype push-to-talk voice-to-text daemon
Documentation=man:voxtype(1)
PartOf=graphical-session.target
After=graphical-session.target pipewire.service pipewire-pulse.service
Wants=ydotool.service

[Service]
# Use Type=notify if voxtype was built with --features systemd
# This allows systemd to wait for the daemon to be fully ready
Type=notify
ExecStart=/usr/bin/voxtype daemon
Restart=on-failure
RestartSec=5
# Timeout for startup notification (model loading can take time)
TimeoutStartSec=60

# Note: User must be in 'input' group for evdev access
# Before enabling this service, run: voxtype setup --download

[Install]
WantedBy=graphical-session.target
